# data_safe_haven/resources/workspace/ansible/files/etc/systemd/system/clamav-clamdscan.service
[Unit]
Description=ClamAV full system scan
Requires=clamav-daemon.service
After=clamav-daemon.service

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/clamdscan --fdpass --multiscan /
